# docker/s6-rc.d/celery/run
#!/usr/bin/with-contenv bash

cd /app/backend

echo "[celery] 等待 FastAPI 服务启动..."
while ! curl -s -f http://localhost:8000/api/v1/status/health > /dev/null 2>&1; do
    echo "[celery] 等待 FastAPI 就绪..."
    sleep 2
done

echo "[celery] FastAPI 已就绪，启动 Celery worker"
exec s6-setuidgid bonita celery -A bonita.worker.celery worker \
    --pool threads \
    --concurrency ${MAX_CONCURRENCY} \
    --events \
    --loglevel DEBUG 